# simeis-server_1.0.0/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable simeis-server.service
systemctl start simeis-server.service
